# resources/saleor/template.service
[Unit]
Description=uWSGI for Saleor API
After=network.target
[Service]
User={USER_NAME}
Restart=always
ExecStart={PYTHON_ENV_PATH}/bin/uwsgi --master --emperor /etc/uwsgi/vassals/prod.ini --uid {USER_NAME} --gid www-data
[Install]
WantedBy=multi-user.target